# files/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/etc/init.d/ssh start

[program:apache2]
command=/etc/init.d/apache2 start

[program:icinga2]
command=/etc/init.d/icinga2 start
